# dev/dam/agent/packaging/dam-agent.service
[Unit]
Description=TooVix DAM Agent (database activity capture)
Documentation=https://toovix.security/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/toovix/agent.env
ExecStart=/usr/local/bin/dam-agent
Restart=always
RestartSec=5

# Runs as root: host (eBPF) mode attaches uprobes to the database's libssl and reads
# /proc/<dbpid>/maps, and network mode opens an AF_PACKET raw socket — both need privilege.
# On the DB host directly there is no container namespace, so no --pid host is required.
User=root
# Explicit bounding set documents exactly what the agent uses (root already holds these; this
# also lets you switch User= to a non-root account by adding AmbientCapabilities=<same set>).
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_NET_RAW CAP_NET_ADMIN
# Light hardening that is still compatible with eBPF + raw sockets.
ProtectHome=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
