# dev/dam/agent/packaging/dam-agent@.service
[Unit]
# Templated so multiple agents (modes) can run on one host without colliding:
#   systemctl enable --now dam-agent@host      -> reads /etc/toovix/agent-host.env
#   systemctl enable --now dam-agent@network   -> reads /etc/toovix/agent-network.env
#   systemctl enable --now dam-agent@proxy     -> reads /etc/toovix/agent-proxy.env
Description=TooVix DAM Agent (%i)
Documentation=https://toovix.security/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/toovix/agent-%i.env
ExecStart=/usr/local/bin/dam-agent
Restart=always
RestartSec=5

# Runs as root: host (eBPF) attaches uprobes to libssl + reads /proc; network opens an
# AF_PACKET raw socket; proxy just listens on a port; audit-forward (AgentLite) reads the
# DB's native audit log, which the DB owns (e.g. /var/log/mysql 750 mysql:adm, files 640).
# CAP_DAC_READ_SEARCH lets root bypass those read/traverse checks (the bounding set caps root
# to exactly this list, so without it root can't read the log — that's a real AgentLite bug).
User=root
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_ADMIN CAP_SYS_PTRACE CAP_NET_RAW CAP_NET_ADMIN CAP_DAC_READ_SEARCH
ProtectHome=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
